# floor796-kiosk.service
[Unit]
Description=Floor796 Kiosk Player
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=HOME=/home/kiosk
Environment=SDL_AUDIODRIVER=dummy
WorkingDirectory=/opt/floor796-kiosk
ExecStart=/opt/floor796-kiosk/run.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
